# packaging/systemd/file-snitch-agent.service.in
[Unit]
Description=File Snitch local agent
Documentation=https://github.com/pkoch/file-snitch/blob/master/docs/services.md
After=default.target

[Service]
Type=simple
ExecStart={{FILE_SNITCH_BIN}} agent --frontend linux-ui
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
